# deploy/single-node/vmon.service
[Unit]
Description=Vibemon Daemon (vmon serve)
After=network.target
Documentation=https://github.com/vibevmm/vibemon

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/var/lib/vmon
EnvironmentFile=-/etc/default/vmon
EnvironmentFile=/etc/vmon/sandbox.env
ExecStart=/usr/local/bin/vmon serve --config /etc/vmon/serve.toml
Restart=always
RestartSec=5

# Privileged requirements for virtualization and networking
LimitNOFILE=65536
LimitNPROC=65536
TasksMax=infinity

# Security/Sandboxing within systemd (restricted but allowing KVM and TAP)
DeviceAllow=/dev/kvm rw
DeviceAllow=/dev/net/tun rw
CapabilityBoundingSet=CAP_NET_ADMIN CAP_SYS_ADMIN CAP_CHOWN CAP_SETUID CAP_SETGID CAP_DAC_OVERRIDE
AmbientCapabilities=CAP_NET_ADMIN CAP_SYS_ADMIN CAP_CHOWN CAP_SETUID CAP_SETGID CAP_DAC_OVERRIDE

[Install]
WantedBy=multi-user.target
